# templates/a36-display.service.template
[Unit]
Description=CPU cooler display controller
After=

[Service]
ExecStart=/usr/bin/a36-display
Restart=always

[Install]
WantedBy=multi-user.target